# etc/systemd/system/test-monitor.service
[Unit]
Description=test processor monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test-monitor.sh
